# unpoweroutter.service
[Unit]
Description=Battery Manager

[Service]
ExecStart=/home/cameron/.cargo/bin/unpoweroutter

[Install]
WantedBy=multi-user.target
